# dist/init/systemd/sigil.service
[Unit]
Description=Sigil web application server

[Service]
User=sigil
Group=sigil

PIDFile=/var/run/sigil/sigil.pid
ExecStart=/usr/bin/sigil

[Install]
WantedBy=multi-user.target
